# launchd/com.codex.slides-relay-agent.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.slides-relay-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd '/ABSOLUTE/PATH/TO/PROJECT' &amp;&amp; ./scripts/slides_relay_agent.sh ./config/relay_agent.env</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/slides-relay-agent.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/slides-relay-agent.err.log</string>
</dict>
</plist>
